# munin/munin.cron
#
# cron-jobs for munin
#

MAILTO=root

*/5 * * * *     munin test -x /usr/bin/munin-cron && /usr/share/munin/munin-update
30 08 * * *     munin test -x /usr/bin/munin-cron && /usr/share/munin/munin-html
